# deploy/systemd/t2sr-web.service
# =====================================================================
# Service systemd — Site T2SR TanStack Start (Bun)
# Copier vers /etc/systemd/system/t2sr-web.service puis:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now t2sr-web
# =====================================================================

[Unit]
Description=T2SR Website (TanStack Start + Bun)
After=network.target

[Service]
Type=simple
User=t2sr
Group=t2sr
WorkingDirectory=/var/www/t2sr

ExecStart=/home/t2sr/.bun/bin/bun run start
Environment=HOST=127.0.0.1
Environment=PORT=4321
EnvironmentFile=-/var/www/t2sr/.env

Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

StandardOutput=journal
StandardError=journal
SyslogIdentifier=t2sr-web

[Install]
WantedBy=multi-user.target
